# supervisord/reboot-api.conf
[program:reboot-api]
autorestart = true
autostart = true
startsecs = 0
directory = /root/reboot-api/
command = /usr/bin/python3 app.py
startretries = 1
environment=PYTHONPATH="%(here)s"
